# dev.conf
[program:dev]
command=/home/ubuntu/dev/venv/bin/uwsgi --ini uwsgi.ini 
directory=/home/ubuntu/dev
;user = forge
stopsignal=QUIT
stdout_logfile=/tmp/dev.stdout.log
stderr_logfile=/tmp/dev.stderr.log
